# specs/hylafax/hylafax_daily.cron
#!/bin/sh

/usr/sbin/faxcron | mail -s "HylaFAX Usage Report" faxmaster

exit 0
